# conf/etc/supervisord.d/cbot.ini
[group:cbot-prod]
programs = cbot-prod
priority = 1

[program:cbot-prod]
process_name = cbot
command = /home/ec2-user/sites/cbot/bin/cbot-server.sh
user = ec2-user
stdout_logfile = /var/log/cbot/cbot.log
stdout_logfile_maxbytes = 0
redirect_stderr = true
priority = 100
